# gnat-env/etc/profile.d/gnat.sh
#!/bin/sh
# setup Ada project path
export GPR_PROJECT_PATH=/usr/share/gpr
export ADA_PROJECT_PATH=/usr/share/gpr
